# dist/lib/systemd/system/opentrickler-led.service
[Unit]
Description=opentrickler leds component
PartOf=opentrickler.service

[Service]
ExecStart=/code/opentrickler/trickler/leds.py /code/opentrickler/opentrickler_config.ini
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
