# dist/lib/systemd/system/opentrickler.service
[Unit]
Description=opentrickler service
After=memcached.service
Requires=memcached.service
Wants=opentrickler-ble.service
Wants=opentrickler-led.service
Wants=chrony.service

[Service]
ExecStart=/code/opentrickler/trickler/main.py /code/opentrickler/opentrickler_config.ini
Restart=always
RestartSec=5

[Install]
Alias=opentrickler
WantedBy=multi-user.target
