# dist/lib/systemd/system/opentrickler-led.service
[Unit]
Description=opentrickler leds component
PartOf=opentrickler.service

[Service]
ExecStart=/code/opentrickler/leds.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
