# dist/lib/systemd/system/opentrickler-ble.service
[Unit]
Description=opentrickler bluetooth component
PartOf=opentrickler.service
After=bluetooth.service
Requires=bluetooth.service
Requires=hciuart.service

[Service]
ExecStart=/code/opentrickler/bluetooth.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
